# rootfs/etc/profile.d/hwclock.sh
# Sync the clock in the Docker Virtual Machine to the system's hardware clock to avoid time drift
hwclock -s
